# config/templates/systemd/backup-date.service.in
[Unit]
Description=Backup the date to a file when we have a valid date time ntp

[Service]
User=root
Group=root
ExecStart=$jaiabot_bin_dir/$exe $args

[Install]
WantedBy=multi-user.target