# config/templates/systemd/goby-app.service.in
[Unit]
Description=$description

BindsTo=jaiabot_health.service
After=jaiabot_health.service
$extra_unit

[Service]
User=$user
Group=$group

EnvironmentFile=$env_file
ExecStartPre=$jaiabot_bin_dir/jaiabot_failure_reporter -C '$gen jaiabot_failure_reporter' --state START --error_code $error_on_fail --service_name $service
ExecStart=$bin_dir/$exe -C '$gen $exe' $extra_flags
ExecStop=$goby_bin_dir/goby_terminate -C '$gen goby_terminate' --target_pid $$MAINPID
ExecStopPost=$jaiabot_bin_dir/jaiabot_failure_reporter -C '$gen jaiabot_failure_reporter' --state STOP --error_code $error_on_fail --service_name $service --service_result $$SERVICE_RESULT
$extra_service

[Install]
WantedBy=$wanted_by
