# config/templates/systemd/goby-app.service.in
[Unit]
Description=$description

BindsTo=gobyd.service
After=gobyd.service
$extra_unit

[Service]
EnvironmentFile=$env_file
ExecStart=$bin_dir/$exe -C '$gen $exe'
ExecStop=$goby_bin_dir/goby_terminate -C '$gen goby_terminate' --target_pid $$MAINPID
$extra_service

[Install]
WantedBy=gobyd.service
