# config/templates/systemd/gobyd.service.in
[Unit]
Description=$description

BindsTo=jaiabot.service
After=jaiabot.service

[Service]
User=$user
Group=$group
Restart=on-failure


EnvironmentFile=$env_file
# Wait 30 seconds before starting
ExecStartPre=/bin/sleep 30
ExecStart=$goby_bin_dir/gobyd -C '$gen gobyd'
ExecStop=$goby_bin_dir/goby_terminate -C '$gen goby_terminate' --target_pid $$MAINPID

[Install]
WantedBy=jaiabot.service
