# config/templates/systemd/gpsd-sim.service.in
[Unit]
Description=$description

BindsTo=jaiabot_health.service
After=jaiabot_health.service
$extra_unit

[Service]
User=$user
Group=$group

EnvironmentFile=$env_file
ExecStart=/usr/sbin/gpsd $$($gen gpsd)

[Install]
WantedBy=jaiabot_health.service
