# config/templates/systemd/gpsd-sim.service.in
[Unit]
Description=$description

[Service]
EnvironmentFile=$env_file
ExecStartPre=/bin/systemctl stop gpsd.socket
ExecStart=/bin/bash -c "exec /usr/sbin/gpsd $$($gen gpsd)"

[Install]
WantedBy=multi-user.target
