# config/templates/systemd/hub-button-trigger.service.in
[Unit]
Description=Hub Button LED Triggers

[Service]
User=root
Group=root

# Run python code within a pre-installed venv
ExecStart=$jaiabot_share_dir/jaiabot/python/venv/bin/python3 $jaiabot_share_dir/jaiabot/python/$subdir/$exe $args

[Install]
WantedBy=multi-user.target
