# config/templates/systemd/jaiabot.service.in
[Unit]
Description=Jaia Base Service

Requires=network-online.target
After=network-online.target

[Service]
User=$user
Group=$group


Type=oneshot
RemainAfterExit=true
ExecStart=/bin/echo "Starting Jaia"

[Install]
WantedBy=multi-user.target
