# config/templates/systemd/jaiabot.service.in
[Unit]
Description=Jaia Base Service

Requires=network-online.target
After=network-online.target

Requires=var-log.mount
After=var-log.mount

Requires=ntp.service
After=ntp.service

[Service]
User=$user
Group=$group


Type=oneshot
RemainAfterExit=true

# Try for 120 seconds for NTP to synchronize
ExecStartPre=/bin/bash -c "for i in {1..12}; do ntpstat && break; echo 'Waiting for NTP sync'; sleep 10; done"

ExecStart=/bin/echo "Starting Jaia"

[Install]
WantedBy=multi-user.target
