# config/templates/systemd/liaison-prelaunch.service.in
[Unit]
Description=$description

Requires=network-online.target
After=network-online.target

Requires=var-log.mount
After=var-log.mount

[Service]
User=$user
Group=$group
Restart=on-failure

EnvironmentFile=$env_file

# Wt appears to not resolve 0::0 until a non-`lo` interface has been assigned an IPv6 address. While this appears to be buggy behavior, the best we can do for now is to ensure that wlan0 is up with IPv6 before we start liaison
ExecStartPre=timeout 60s bash -c 'until ip address show dev wlan0 | grep "inet6"; do sleep 1; done'

ExecStart=$bin_dir/$exe -C '$gen goby_liaison_prelaunch' $extra_flags
$extra_service
WorkingDirectory=$ansible_dir


[Install]
WantedBy=multi-user.target
