# config/templates/systemd/moosdb-sim.service.in
[Unit]
Description=$description

BindsTo=jaiabot_health.service
After=jaiabot_health.service


[Service]
User=$user
Group=$group

EnvironmentFile=$env_file

# Wait 5 seconds before starting
ExecStartPre=/bin/sleep 5
ExecStartPre=/bin/bash -c '$gen moos_sim > $moos_sim_file'

ExecStartPre=$jaiabot_bin_dir/jaiabot_failure_reporter -C '$gen jaiabot_failure_reporter' --state START --error_code $error_on_fail --service_name $service
ExecStart=$moos_bin_dir/MOOSDB $moos_sim_file
ExecStopPost=$jaiabot_bin_dir/jaiabot_failure_reporter -C '$gen jaiabot_failure_reporter' --state STOP --error_code $error_on_fail --service_name $service --service_result $$SERVICE_RESULT
TimeoutSec=180

[Install]
WantedBy=jaiabot_health.service
