# config/templates/systemd/moosdb.service.in
[Unit]
Description=$description

BindsTo=jaiabot.service
After=jaiabot.service

[Service]
EnvironmentFile=$env_file
ExecStartPre=/bin/bash -c '$gen moos > $moos_file'
ExecStartPre=/bin/bash -c '$gen bhv > $bhv_file'
ExecStart=$moos_bin_dir/MOOSDB $moos_file

[Install]
WantedBy=jaiabot.service
