# config/templates/systemd/py-app.service.in
[Unit]
Description=$description

BindsTo=jaiabot_gobyd.service
After=jaiabot_gobyd.service
$extra_unit

[Service]
EnvironmentFile=$env_file
# Run python code within a pre-installed venv
ExecStart=$jaiabot_share_dir/jaiabot/python/venv/bin/python3 $jaiabot_share_dir/jaiabot/python/$subdir/$exe $args
$extra_service

[Install]
WantedBy=jaiabot_gobyd.service
