# config/templates/systemd/py-app.service.in
[Unit]
Description=$description

BindsTo=jaiabot_gobyd.service
After=jaiabot_gobyd.service
$extra_unit

[Service]
EnvironmentFile=$env_file
ExecStart=$jaiabot_bin_dir/jaia_py/$subdir/$exe $args
$extra_service

[Install]
WantedBy=jaiabot_gobyd.service
